# DIST/serverchecker-0.1/debian/serverchecker/etc/init.d/serverchecker
#!/bin/sh

### BEGIN INIT INFO

# Provides:          ServerChecker

# Required-Start:    $local_fs $network $named $time $syslog

# Required-Stop:     $local_fs $network $named $time $syslog

# Default-Start:     2 3 4 5

# Default-Stop:      0 1 6

# Description:       ServerChecker application with embedded web service.

### END INIT INFO



FOLDER=/usr/bin/ServerChecker
RUNAS=root
SERVICE_NAME=serverchecker
PATH_TO_JAR=/usr/bin/ServerChecker/ServerChecker.jar
PID_PATH_NAME=/var/run/serverchecker.pid
LOGFILE=/var/log/SERVERCHECKER.log




PID_PATH_NAME=/tmp/MyService-pid
case $1 in
    start)
        echo "Starting $SERVICE_NAME ..."
        if [ ! -f $PID_PATH_NAME ]; then
	    cd $FOLDER
            nohup java -jar $PATH_TO_JAR >> $LOGFILE 2>&1 &
                        echo $! > $PID_PATH_NAME
            echo "$SERVICE_NAME started ..."
        else
            echo "$SERVICE_NAME is already running ..."
        fi
    ;;
    stop)
        if [ -f $PID_PATH_NAME ]; then
            PID=$(cat $PID_PATH_NAME);
            echo "$SERVICE_NAME stoping ..."
            kill $PID;
            echo "$SERVICE_NAME stopped ..."
            rm $PID_PATH_NAME
        else
            echo "$SERVICE_NAME is not running ..."
        fi
    ;;
    restart)
        if [ -f $PID_PATH_NAME ]; then
            PID=$(cat $PID_PATH_NAME);
            echo "$SERVICE_NAME stopping ...";
            kill $PID;
            echo "$SERVICE_NAME stopped ...";
            rm $PID_PATH_NAME
            echo "$SERVICE_NAME starting ..."
            nohup java -jar $PATH_TO_JAR /tmp 2>> /dev/null >> /dev/null &
                        echo $! > $PID_PATH_NAME
            echo "$SERVICE_NAME started ..."
        else
            echo "$SERVICE_NAME is not running ..."
        fi
    ;;
esac

